# etc/systemd/system/gpuboot.service
[Unit]
Description=GPUBOOT: apply power saving features . . .
After=basic.target

[Service]
User=root
Type=simple
ExecStart=/bin/gpuboot

[Install]
WantedBy=basic.target
